# systemd/tg_sticket_bot.service
[Unit]
Description=Telegram sticker bot

[Service]
Type=simple
ExecStart=/opt/tg_sticker_bot/start.sh

[Install]
WantedBy=multi-user.target
